# Intentional/Library/LaunchAgents/com.intentional.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.intentional.watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>FLAG="$HOME/Library/Application Support/Intentional/strict-mode"; NORELAUNCH="${TMPDIR}intentional-no-relaunch"; [ -f "$FLAG" ] &amp;&amp; ! [ -f "$NORELAUNCH" ] &amp;&amp; ! pgrep -x Intentional &gt;/dev/null &amp;&amp; open -a Intentional</string>
	</array>
	<key>StartInterval</key>
	<integer>10</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/intentional-watchdog.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/intentional-watchdog.log</string>
</dict>
</plist>
